# mw.service
[Unit]
Description=Gmail synchronization service

[Service]
Type=oneshot
ExecStart=/usr/bin/mbsync -Va --push --create
